# frozen-pipe-monitor.service
[Unit]
Description=Frozen Pipe Weather Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/frozen
ExecStart=/usr/bin/python3 /home/pi/frozen/main.py monitor --interval 30
Restart=always
RestartSec=60
StandardOutput=append:/home/pi/frozen/monitor.log
StandardError=append:/home/pi/frozen/monitor.log

# Resource limits (optional, for extra safety)
MemoryLimit=100M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
